# bme280.service
[Unit]
Description=PolandAOD power on BME280

[Service]
Type=simple
ExecStart=/usr/local/bin/get_bme280.py &

[Install]
WantedBy=multi-user.target

