# opc.service
[Unit]
Description=PolandAOD powers on OPC

[Service]
Type=simple
ExecStart=/usr/local/bin/get_opc.py &

[Install]
WantedBy=multi-user.target

