# pinger.service
[Unit]
Description=PolandAOD pinger 

[Service]
Type=simple
ExecStart=/usr/local/bin/pinger.sh &

[Install]
WantedBy=multi-user.target

