# rtc.service
[Unit]
Description=RTC with battery backup

[Service]
Type=simple
ExecStart=/usr/local/bin/rtc.sh &

[Install]
WantedBy=multi-user.target

